# deploy/instagram-publish-scheduler.service
[Unit]
Description=Instagram Publish Scheduler - process scheduled posts
After=instagram-logger.service

[Service]
Type=oneshot
EnvironmentFile=/home/christof21/instagram-post-refiner/.env.local
ExecStart=/usr/bin/curl -sf -H "x-cron-secret: ${CRON_SECRET}" http://localhost:3000/api/cron/publish
TimeoutStartSec=55
