# deploy/instagram-snapshot.service
[Unit]
Description=Instagram Post Logger — Nightly account snapshot
After=network-online.target instagram-logger.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/home/christof21/instagram-post-refiner/.env.local
ExecStart=/usr/bin/curl -sf -X POST -H "x-cron-secret: ${CRON_SECRET}" http://localhost:3000/api/instagram/snapshot
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
