# deploy/instagram-poll-new-posts.service
[Unit]
Description=Instagram Post Logger — Poll for new posts
After=network-online.target instagram-logger.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/home/christof21/instagram-post-refiner/.env.local
ExecStart=/usr/bin/curl -sf -H "x-cron-secret: ${CRON_SECRET}" http://localhost:3000/api/cron/poll-new-posts
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
