# deploy/instagram-metrics-sync.service
[Unit]
Description=Instagram Post Logger — Nightly cron (snapshot + backfill + metrics)
After=network-online.target instagram-logger.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -sf http://localhost:3000/api/cron/nightly
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
