# deploy/instagram-snapshot.service
[Unit]
Description=Instagram Post Logger — Nightly account snapshot
After=network-online.target instagram-logger.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s -X POST http://localhost:3000/api/instagram/snapshot
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
