# deploy/instagram-publish-scheduler.timer
[Unit]
Description=Instagram Publish Scheduler (every minute)

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
